# deploy/oko-bot.service
[Unit]
Description=Oko Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/oko
EnvironmentFile=/opt/oko/.env
ExecStart=/opt/oko/venv/bin/python -m src.bot.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oko-bot

[Install]
WantedBy=multi-user.target
